# goterra.service.example
[Unit]
Description=goterra

[Service]
Type=simple
Restart=always
RestartSec=5s
WorkingDirectory=path_to_goterra
ExecStart=path_to_goterra/goterra

[Install]
WantedBy=multi-user.target

